# deploy/stats-consolidation.cron
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

0 3 * * * root stats_consolidation --config_file /etc/stats-consolidation.conf > $BASE_PATH/log/out.txt

